# services/rsscrawler/run
#!/bin/bash
exec /sbin/setuser abc /usr/bin/python /config/RSScrawler.py --log-level=$LOGLEVEL
