# root/etc/services.d/rsscrawler/run
#!/usr/bin/with-contenv bash

if [[ -v LOGLEVEL ]]; then
  exec s6-setuidgid abc python3 /config/RSScrawler.py --docker --log-level=$LOGLEVEL $PARAMETER
else
  exec s6-setuidgid abc python3 /config/RSScrawler.py --docker $PARAMETER
fi
