# root/etc/services.d/rsscrawler/run
#!/usr/bin/with-contenv bash

if [[ -v LOGLEVEL ]]; then
  exec rsscrawler --docker --log-level=$LOGLEVEL $PARAMETER
else
  exec rsscrawler --docker $PARAMETER
fi
